# imagefiles-supervisord/jupyter.conf
[program:jupyter]
command=jupyter notebook
