# imagefiles-supervisord/jupyter.conf
[program:jupyter]
user=x2go
directory=/home/x2go
environment=HOME="/home/x2go",USER="x2go",jupyter_config="/home/x2go/.jupyter"
command=jupyter notebook
